# roles/nginx/files/etc-nginx-session_tickets_keys.automount
[Unit]
Description=session ticket keys automount

[Automount]
Where=/etc/nginx/session_tickets_keys

[Install]
WantedBy=multi-user.target
